# etc/rc.local
#!/bin/sh

iptables-restore < /etc/iptables.sav

exit 0